# ops/systemd/internal-image-studio-backup.service
[Unit]
Description=Internal Image Studio backup
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/internal-image-studio/backup.env
ExecStart=/opt/internal-image-studio/ops/backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
