# ops/systemd/internal-image-studio-restore-check.service
[Unit]
Description=Cornfield monthly backup restore verification
After=docker.service network-online.target internal-image-studio-maintenance-recovery.service
Requires=docker.service
Wants=internal-image-studio-maintenance-recovery.service

[Service]
Type=oneshot
EnvironmentFile=/etc/internal-image-studio/backup.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
StateDirectory=internal-image-studio
StateDirectoryMode=0700
ExecStart=/opt/internal-image-studio/ops/restore-check.sh
ExecStopPost=/opt/internal-image-studio/ops/restore-check-recovery.sh
TimeoutStartSec=8h
TimeoutStopSec=8m
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
