# ops/systemd/internal-image-studio-backup.timer
[Unit]
Description=Daily Internal Image Studio backup

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
